# release/arch/nvcontrol-monitor.service
[Unit]
Description=nvcontrol GPU Monitoring Service
Documentation=https://github.com/GhostKellz/nvcontrol
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/nvctl monitor start
Restart=on-failure
RestartSec=10s

# Resource limits
Nice=-10
IOSchedulingClass=best-effort
IOSchedulingPriority=2

# Security
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nvcontrol-monitor

[Install]
WantedBy=graphical-session.target
